# debian/mydlp.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

25 6	* * *	root	/usr/sbin/mydlp-index


